# contrib/nexus-open.service
[Unit]
Description=Nexus Open — iCUE Nexus companion daemon
Documentation=https://github.com/mantonx/nexus-next
After=graphical-session.target
Wants=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.local/bin/nexus-open --tray
Restart=on-failure
RestartSec=3s

[Install]
WantedBy=graphical-session.target
