# packaging/deb/DEBIAN/postinst
#!/bin/bash
set -e

echo "Setting up Nexus Open..."

# Ensure plugdev group exists — fallback for headless installs where
# TAG+="uaccess"/logind is not active.
if ! getent group plugdev > /dev/null 2>&1; then
    echo "Creating plugdev group..."
    groupadd -r plugdev
fi

# Reload udev so the new rules take effect immediately.
# Guard against container/CI environments where udevd is not running.
if command -v udevadm > /dev/null 2>&1; then
    echo "Reloading udev rules..."
    udevadm control --reload-rules 2>/dev/null || true
    udevadm trigger --subsystem-match=usb 2>/dev/null || true
    udevadm trigger --subsystem-match=hidraw 2>/dev/null || true
fi

# Update desktop and icon caches if available.
command -v update-desktop-database > /dev/null 2>&1 && \
    update-desktop-database -q /usr/share/applications || true
command -v gtk-update-icon-cache > /dev/null 2>&1 && \
    gtk-update-icon-cache -q -t -f /usr/share/icons/hicolor || true

echo ""
echo "Nexus Open installed successfully!"
echo ""
echo "Unplug and replug your iCUE Nexus — it will be accessible immediately"
echo "in your desktop session (TAG+=\"uaccess\" handles this automatically)."
echo ""
echo "On headless or multi-user systems, also run:"
echo "  sudo usermod -a -G plugdev \$USER"
echo "and log out/in for the group change to take effect."
echo ""
echo "To start Nexus Open at login:"
echo "  systemctl --user enable --now nexus-open"
echo ""

exit 0
