# packaging/udev/99-corsair-nexus.rules
# udev rules for Corsair iCUE Nexus (VID: 0x1b1c, PID: 0x1b8e)
#
# TAG+="uaccess" grants access to the active desktop session automatically
# via logind — this is the primary access mechanism on all systemd desktops.
# No group membership is required for normal desktop use.
#
# MODE="0600" is a root-only fallback: it prevents world-read while still
# allowing logind/uaccess to set per-session ACLs.  On desktops, uaccess
# is sufficient; users should NOT need to join plugdev.
#
# GROUP="plugdev" + MODE="0660" is only needed on headless or multi-user
# setups where no logind session is active (e.g. a server running the daemon
# under a dedicated service account).  Enable it manually in that case.
#
# RUN+="usbreset" issues a USB reset on each attach so any stale usbfs handle
# from a previous unclean exit is cleared and hid-generic rebinds cleanly.
# Requires the usbreset binary (shipped in the nexus-open package).

SUBSYSTEM=="usb", ATTRS{idVendor}=="1b1c", ATTRS{idProduct}=="1b8e", \
    MODE="0600", TAG+="uaccess", \
    RUN+="/usr/bin/usbreset /dev/bus/usb/$env{BUSNUM}/$env{DEVNUM}"

# hidraw nodes do not always inherit ACLs from the parent USB device
# (distro-dependent). This rule ensures direct hidraw access works everywhere.
SUBSYSTEM=="hidraw", ATTRS{idVendor}=="1b1c", ATTRS{idProduct}=="1b8e", \
    MODE="0600", TAG+="uaccess"
