# packaging/udev/99-corsair-nexus.rules
# udev rules for Corsair iCUE Nexus (VID: 0x1b1c, PID: 0x1b8e)
#
# TAG+="uaccess" grants access to the active desktop session automatically
# via logind — no group membership required on any systemd distro.
#
# GROUP="plugdev" is a fallback for headless / multi-user setups where no
# logind session is present. The group is created by the package installer
# if it doesn't already exist.
#
# RUN+="usbreset" issues a USB reset on each attach so any stale usbfs handle
# from a previous unclean exit is cleared and hid-generic rebinds cleanly.
# This eliminates the need to physically replug the device after the app
# restarts (e.g. after a crash or SIGKILL).

SUBSYSTEM=="usb", ATTRS{idVendor}=="1b1c", ATTRS{idProduct}=="1b8e", \
    MODE="0660", GROUP="plugdev", TAG+="uaccess", \
    RUN+="/usr/bin/usbreset /dev/bus/usb/$env{BUSNUM}/$env{DEVNUM}"

# hidraw nodes do not always inherit ACLs from the parent USB device
# (distro-dependent). This rule ensures direct hidraw access works everywhere.
SUBSYSTEM=="hidraw", ATTRS{idVendor}=="1b1c", ATTRS{idProduct}=="1b8e", \
    MODE="0660", GROUP="plugdev", TAG+="uaccess"
